# deploy/oracle/bithumb-tradbot.service.example
[Unit]
Description=Bithumb Tradbot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bithumb_tradbot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/bithumb_tradbot/.venv/bin/python /home/ubuntu/bithumb_tradbot/main.py --config /home/ubuntu/bithumb_tradbot/config.yaml
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
